# src/opt/switchportmapper/etc/switchportmapper.cron
*/15 * * * * nobody /opt/switchportmapper/bin/doSwitchMap 
